# webapp/supervisor.conf
[program:suwa]
command=node server.js
directory=/home/pi/suwa/webapp/api
autostart=true
autorestart=true

